# columnstore_base/service/columnstore/finish
#!/bin/sh
rm -f /var/lock/subsys/columnstore

